# scripts/supervisor.conf
[program:sws-server]
command=/usr/local/bin/sws-server
directory=/tmp
user=root
autostart=true
autorestart=true