# scripts/supervisor.conf
[program:shadowx-server]
command=/usr/local/bin/shadowx-server
directory=/tmp
user=root
autostart=true
autorestart=true