# scripts/supervisor.conf
[program:ssx-server]
command=/usr/local/bin/ssx-server
directory=/tmp
user=root
autostart=true
autorestart=true